# imageroot/systemd/user/phpFpm.service
[Unit]
Description=Timer of the update of phpFpm pool

[Service]
Type=oneshot
ExecStart=runagent update-php-fpm
SyslogIdentifier=%u/%N
