# imageroot/systemd/user/phpFpm.timer
[Unit]
Description=Timer of the update of phpFpm pool

[Timer]
OnCalendar=monthly
RandomizedDelaySec=3600
[Install]
WantedBy=timers.target
